# radio_track/submodules/nprotocol/debian/noderesolver.prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade)

    systemctl stop noderesolver || true
    systemctl disable noderesolver || true

    ;;
esac
